# cron/staging/www-data.cron
*/5 * * * * export PYTHONPATH=/deploy/staging/; export DEPLOYMENT=staging; /deploy/staging/vib/monitoring/MonitorQueues-wrapper.py >> /tmp/MonitorQueues-staging-$(/bin/date +"\%F").log 2>&1
