# cron/www-data.prod.cron
0 */2 * * * export PYTHONPATH=/deploy/prod/; /deploy/prod/vib/utils/notify_mturk_workers.py >> /tmp/notify_mturk_workers-$(/bin/date +"\%F").log 2>&1
