# cron/www-data.local.cron
0 */2 * * * export PYTHONPATH=/deploy/local/; /deploy/local/vib/utils/notify_mturk_workers.py >> /tmp/notify_mturk_workers-$(/bin/date +"\%F").log 2>&1
